# rak/rak-cellular-utilities/files/uci-defaults/99-rak-cellular-rpcd
#!/bin/sh
# Ce script redémarre rpcd puis s'auto-détruit (mécanique uci-defaults)
/etc/init.d/rpcd restart
exit 0